# platform/systemd.services/localstack.service
[Unit]
Description=localstack

[Service]
ExecStart=/home/docker/.local/bin/localstack start
User=docker

[Install]
WantedBy=multi-user.target